# scripts/crontab
# Run secret rotation every 30 seconds
* * * * * /scripts/rotate-secrets.sh
* * * * * sleep 30 && /scripts/rotate-secrets.sh 